# puppet/libvirt/templates/libvirt-guests.erb
# chkconfig: 345 99 11
# description: Priority level for libvirt-guests service.
