# puppet/libvirt/templates/libvirtd.erb
# chkconfig: 345 97 12
# description: The priority level for libvirtd service.
